# plists/me.damca.brewd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>KeepAlive</key>
	<dict>
		<key>SuccessfulExit</key>
		<false/>
	</dict>
	<key>Label</key>
	<string>me.damca.brewd</string>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/bin/brew</string>
		<string>update</string>
	</array>
	<key>StandardErrorPath</key>
	<string>/usr/local/var/brewd/error.log</string>
	<key>StandardOutPath</key>
	<string>/usr/local/var/brewd/output.log</string>
	<key>StartCalendarInterval</key>
	<dict>
		<key>Hour</key>
		<integer>10</integer>
		<key>Minute</key>
		<integer>0</integer>
		<key>Weekday</key>
		<integer>2</integer>
	</dict>
	<key>ThrottleInterval</key>
	<integer>60</integer>
</dict>
</plist>
